# init/init.upstart
description     "Barsystem"
author          "Jasper Seidel <jawsper@jawsper.nl>"

# no start option as you might not want it to auto-start
# This might not be supported - you might need a: start on runlevel [3]
start on runlevel [2345]

# if you want it to automatically restart if it crashes, leave the next line in
respawn

script
    su -c {{project_directory}}/launcher.py tkkrlab
end script